# root/projectroot/etc/systemd/system/systemd-journald.service.d/persistent-journal.conf
[Unit]
RequiresMountsFor=/var/log/journal
