# root/base/projectroot/usr/lib/systemd/system/app.socket
[Unit]
PartOf=app.service
RefuseManualStart=true
RefuseManualStop=true

[Socket]
ListenStream=33000
BindToDevice=sys
