# root/base/projectroot/usr/lib/systemd/system/codesys.service

[Service]
WorkingDirectory=/var/codesys
ExecStart=/opt/codesys/bin/codesyscontrol.bin --debug --printlostcycles /etc/CODESYSControl.cfg
