# root/base/projectroot/usr/lib/systemd/system/debug.target
[Unit]
Wants=studio.socket

[Install]
WantedBy=multi-user.target
