# root/base/projectroot/usr/lib/systemd/system/gui-admin.socket
[Unit]
StopWhenUnneeded=true

[Socket]
ListenStream=8080
BindToDevice=sys
