# root/base/projectroot/usr/lib/systemd/system/gui-admin.service
[Unit]
PartOf=gui-admin.socket

[Service]
SyslogIdentifier=gui-admin
WorkingDirectory=/usr/lib/gui
ExecStart=/usr/bin/python -Bu /usr/lib/gui admin
