# root/base/projectroot/usr/lib/systemd/system/gui-hmi.service
[Unit]
PartOf=gui-hmi.socket

[Service]
SyslogIdentifier=gui-hmi
WorkingDirectory=/usr/lib/gui
ExecStart=/usr/bin/python -Bu /usr/lib/gui hmi
