# root/base/projectroot/usr/lib/systemd/system/rt-setup-irq@.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/rt-setup-irq %i
StandardOutput=null
