# root/base/projectroot/usr/lib/systemd/system/studio.service
[Unit]
PartOf=studio.socket

[Service]
WorkingDirectory=/usr/lib/gui
ExecStart=/usr/lib/gui/studio-server
