# root/base/projectroot/usr/lib/systemd/system/user@admin.target
[Unit]
StopWhenUnneeded=true

Wants=gui-admin.socket
