# root/base/projectroot/usr/lib/systemd/system/user@oem.target
[Unit]
StopWhenUnneeded=true

Wants=gui-studio.socket

Wants=user@admin.target
