# root/projectroot/usr/lib/systemd/system/arduino-update.socket
[Socket]
BindToDevice=sys
ListenStream=8101
Accept=yes
