# root/projectroot/usr/lib/systemd/system/arduino-update@.service
[Service]
WorkingDirectory=/usr/lib/arduino
ExecStart=/usr/lib/arduino/arduino-update
