# root/projectroot/usr/lib/systemd/system/codesys-log.service
[Service]
Type=notify
RuntimeDirectory=codesys-log
ExecStart=/opt/codesys/bin/log-journal.py
