# root/projectroot/usr/lib/systemd/system/codesys-proxy.service
[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd --connections-max=5 --exit-idle-time=5 127.0.0.1:11749
