# root/projectroot/usr/lib/systemd/system/codesys-proxy.socket
[Unit]
StopWhenUnneeded=true

[Socket]
ListenStream=11740
BindToDevice=sys
