# root/projectroot/usr/lib/systemd/system/esp32-ota.socket
[Socket]
BindToDevice=sys
ListenStream=8102
Accept=yes
