# root/projectroot/usr/lib/systemd/system/esp32-ota@.service
[Service]
ExecStart=/usr/lib/esp32/esp32-ota
