# root/projectroot/usr/lib/systemd/system/etc.mount
[Unit]
RequiresMountsFor=/var/etc

[Mount]
What=overlay
Where=/etc
