# root/projectroot/usr/lib/systemd/system/gui.socket
[Unit]
ConditionPathExists=/usr/lib/gui/gui

[Socket]
ListenStream=80
BindToDevice=sys

[Install]
WantedBy=sockets.target
