# root/projectroot/usr/lib/systemd/system/gui.service
[Unit]
PartOf=gui.socket

[Service]
WorkingDirectory=/usr/lib/gui
ExecStart=/usr/lib/gui/gui
