# root/projectroot/usr/lib/systemd/system/mc-state.service
[Unit]
Conflicts=mc-state.socket

[Service]
ExecStart=/usr/sbin/mc-state mc 33000
ExecStopPost=-systemctl start --no-block mc-state.socket
