# root/projectroot/usr/lib/systemd/system/mc-shaker.service
[Unit]
StopWhenUnneeded=true

[Service]
ExecStart=/usr/sbin/mc-shaker
