# root/projectroot/usr/lib/systemd/system/mc-update.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/mc-update
