# root/projectroot/usr/lib/systemd/system/reboot.service
[Service]
Type=oneshot
ExecStart=/bin/systemctl --no-block reboot
