# root/projectroot/usr/lib/systemd/system/shaker.service
[Unit]
StopWhenUnneeded=true

[Service]
ExecStart=/usr/sbin/shaker
