# root/projectroot/usr/lib/systemd/system/studio.socket
[Unit]
StopWhenUnneeded=true

[Socket]
ListenStream=8000
BindToDevice=sys
