# boot/projectroot/etc/inittab

# mount kernel filesystems
::sysinit:/bin/mount /dev
::sysinit:/bin/mount /proc
::sysinit:/bin/mount /sys

# run rcS
::once:/bin/cttyhack /etc/init.d/rcS
